# system/Library/LaunchDaemons/com.yodelcms.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.yodelcms.server</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/local/bin/yodel_command_runner</string>
		<string>server</string>
	</array>
	<key>KeepAlive</key>
	<true/>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardOutPath</key>
	<string>/var/log/yodel.log</string>
	<key>StandardErrorPath</key>
	<string>/var/log/yodel.log</string>
	<dict>
	  <key>INLINEDIR</key>
	  <string><%= sites_root %></string>
	</dict>
</dict>
</plist>
